# systemd/agent-bridge-gemini.service
[Unit]
Description=Telegram bridge for Gemini
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
EnvironmentFile=/etc/default/agent-bridge-gemini
ExecStart=/usr/bin/env bash -lc 'cd "${BRIDGE_PROJECT_DIR:?}" && exec node src/index.js'
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
