# systemd/agent-bridge-worker-bot.service
[Unit]
Description=Telegram bridge — autonomous worker bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=BRIDGE_USER
EnvironmentFile=-/etc/default/agent-bridge-shared
EnvironmentFile=/etc/default/agent-bridge-worker-bot
ExecStart=/usr/bin/env bash -lc 'export PATH="$(dirname "${NODE_BIN:?}"):$PATH"; cd "${BRIDGE_PROJECT_DIR:?}" && exec "${NODE_BIN:?}" ./node_modules/tsx/dist/cli.mjs src/index-worker.ts'
Restart=always
RestartSec=5
KillMode=control-group
Environment=BRIDGE_WORKSPACE_LOCK_MODE=on
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
