# systemd/agent-bridge-interactive.service
[Unit]
Description=Telegram bridge — unified interactive bot (switchable CLI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=BRIDGE_USER
EnvironmentFile=-/etc/default/agent-bridge-shared
EnvironmentFile=/etc/default/agent-bridge-interactive
ExecStart=/usr/bin/env bash -lc 'export PATH="$(dirname "${NODE_BIN:?}"):$PATH"; cd "${BRIDGE_PROJECT_DIR:?}" && exec "${NODE_BIN:?}" ./node_modules/tsx/dist/cli.mjs src/index-interactive.ts'
Restart=always
RestartSec=5
KillMode=control-group
Environment=BRIDGE_WORKSPACE_LOCK_MODE=off
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
