# systemd/agent-bridge-gemini.service
[Unit]
Description=Telegram bridge for Gemini
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Environment=BRIDGE_ROOT_DIR=/home/openclaw
Environment=BRIDGE_PROJECT_DIR=/home/openclaw/.openclaw/workspace/projects/agent-bridge
EnvironmentFile=/home/openclaw/.openclaw/workspace/projects/agent-bridge/.env.gemini
WorkingDirectory=/home/openclaw/.openclaw/workspace/projects/agent-bridge
ExecStart=/usr/bin/env node src/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
